# dist/etc/skel/.config/autostart/Picom.desktop
[Desktop Entry]
Exec=picom --legacy-backends --backend=xr_glx_hybrid
X-XFCE-Autostart-Override=true
Hidden=false
